# share/etc.systemd.system.hypno.service
[Unit]
Description=Hypno, a Lua-based web and application server.
DefaultDependencies=no
After=network.target

[Service]
Type=simple
ExecStart=@@PREFIX@@/bin/hypno-server --start --port 80 --user @@SRVUSER@@ --group @@SRVGROUP@@ --config /var/www/config.lua
Restart=on-failure

[Install]
WantedBy=multi-user.target
